# etc/rc.local
ip rule add fwmark 0x01/0x01 table 100
ip route add local 0.0.0.0/0 dev lo table 100
iptables-restore < /etc/shadowsocks-libev/iptables.ss-redir
